# service/vllm-glm-4-flash.service
[Unit]
Description=vLLM GLM-4.7-Flash Service
After=docker.service
Requires=docker.service
Conflicts=vllm-qwen3-coder-next.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=jack

ExecStart=/home/jack/deploy/vllm/spark-vllm-docker/vllm-start-glm-4-flash.sh
ExecStop=/usr/bin/docker stop --time=60 vllm_node
ExecStopPost=/usr/bin/docker rm -f vllm_node

TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
